# etc/supervisor/conf.d/camera.conf
[program:camera_streamer]
command=python /home/pi/utils/flask-video-streaming/app.py
user=root
directory=/home/pi
autostart=true
environment=PYTHONUNBUFFURED=1
